# service.conf
[program:goberry]
command=/%(ENV_GOPATH)s/bin/goberry
stdout_events_enabled = true
stderr_events_enabled = true
